# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=gunicorn -c gunicorn_conf.py

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"